# pbx-startup-tests.service
[Unit]
Description=Warden Voip System Startup Tests
After=network.target

[Service]
Type=oneshot
# IMPORTANT: Update these paths to match your PBX installation directory
# Common locations: /opt/pbx, /usr/local/pbx, /home/pbx/PBX
WorkingDirectory=/opt/pbx
ExecStart=/opt/pbx/run_startup_tests.sh
# IMPORTANT: Update this to the user who owns the PBX files
# Common users: pbx, root, your-username
User=pbx
StandardOutput=journal
StandardError=journal
# Don't fail the boot if tests fail (exit codes 0 and 1 are both success)
SuccessExitStatus=0 1

[Install]
WantedBy=multi-user.target
